# back/.docker/crons/queue_all_stats.cron
*/5 * * * * /usr/local/bin/php /var/www/bin/console collect:stats:detailed
