# back/.docker/crons/queue_server_stats.cron
*/5 * * * * /usr/local/bin/php /var/www/bin/console queue:produce:server_stats
